# llama-server.service
# llama-server.service - Updated with 131k context window

[Unit]
Description=llama.cpp OpenAI-compatible server
After=network.target

[Service]
Type=simple
User=localadmin
Group=localadmin
ExecStart=/opt/llama.cpp/build/bin/llama-server -m /opt/models/gguf/qwen3.5:latest.gguf --mmproj /opt/models/gguf/mmproj-F16.gguf --host 172.16.254.100 --port 11434 --ctx-size 131072 --context-shift --slot-idle-timeout 600 --slot-prefill-max 4 --slot-max-queue 10
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target