# wiki/wiki.service
[Unit]
Description=Workspace Wiki (Next.js)
After=network.target

[Service]
Type=simple
User=localadmin
WorkingDirectory=/home/localadmin/.openclaw/workspace/wiki
ExecStart=/usr/bin/node /home/localadmin/.openclaw/workspace/wiki/node_modules/.bin/next start -p 3000
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
