# projects/realtime-dashboard/realtime-dashboard.service
[Unit]
Description=Real-time Systems Dashboard Server
After=network.target

[Service]
Type=simple
User=localadmin
WorkingDirectory=/home/localadmin/.openclaw/workspace/skills/daily-systems-status
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
